# venn-api.service
[Unit]
Description=Venn API (FastAPI/uvicorn)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/venn/backend
Environment=SECRET_KEY=cacdec6c5f54d602ecd003355d3bfc9725b710b2988660a974beb45dfaf4b134
ExecStart=/usr/bin/python3 -m uvicorn main:app --host 127.0.0.1 --port 8000
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
